# phpvirtualbox/rootfs/etc/services.d/nginx/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: phpVirtualBox
# Runs the Nginx daemon
# ==============================================================================

# Wait for PHP-FPM to become available
bashio::net.wait_for 9001

bashio::log.info "Starting Nginx..."

exec nginx
